# src/sd/firmware/etc/profile.d/config.sh
set -o allexport
[ -f $SD/config.cfg ] && source $SD/config.cfg
set +o allexport
