# src/sd/firmware/etc/profile.d/path.sh
export PATH=$SD/firmware/bin:$PATH
export LD_LIBRARY_PATH=$SD/firmware/lib:$LD_LIBRARY_PATH
